# backend/llamacpp-server.service
[Unit]
Description=llama.cpp Server for JARVIS
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/llama.cpp
ExecStart=/home/pi/llama.cpp/server \
    -m /home/pi/llama.cpp/models/qwen-1_5b-chat-q4_0.gguf \
    -c 2048 \
    -ngl 0 \
    --host 0.0.0.0 \
    --port 8080 \
    -t 4
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
